# debian/nova-xvpvncproxy.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d nova-xvpvncproxy remove
fi

#DEBHELPER#
