# debian/nova-volume.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d nova-volume remove
fi

#DEBHELPER#
